# meta-intel-edison-distro/recipes-connectivity/mostfun-panel/files/panel-install.service
[Unit]
Description=panel install service
After=TFT_LCD.service

[Service]
Type=oneshot
ExecStart=/etc/panel-install.sh
ExecStartPost=/bin/systemctl disable panel-install.service
StandardOutput=journal+console

[Install]
WantedBy=basic.target
