# meta-intel-edison-distro/recipes-core/base-files/base-files/media-sdcard.automount
[Unit]
Description=Automount for USB storage
ConditionPathExists=/dev/sdb1
RefuseManualStart=true
RefuseManualStop=true
After=local-fs.target

[Automount]
Where=/media/usb

[Install]
WantedBy=local-fs.target
