# meta-intel-edison-distro/recipes-core/base-files/base-files/media-sdcard.mount
[Unit]
Description=Mount for SDCard
BindsTo=dev-sdb1.device
Conflicts=media-sdcard.automount

[Mount]
What=/dev/sdb1
Where=/media/usb
TimeoutSec=3
